# tools/telegram-bot/claude-tg-bot.service
[Unit]
Description=Claude Code Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/root/telegram-claude
EnvironmentFile=/root/telegram-claude/.env
ExecStart=/usr/bin/python3 /root/telegram-claude/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
